# etc/rc.local

#sudo /usr/local/bin/pip install -U setuptools
#sudo /usr/local/bin/easy_install awscli docker boto3
#sudo yum install -y git systat screen

ZONE=$(curl -sS http://169.254.169.254/latest/meta-data/placement/availability-zone)
export AWS_DEFAULT_REGION=$(echo $ZONE|sed 's/.$//')

DIR=/home/ec2-user/ecs-ebs-manager
git -C $DIR fetch
git -C $DIR pull
$DIR/run-docker-ebs-manager.py &

